# deploy/kokovoicebot.service
[Unit]
Description=Kokovoicebot Oracle Control Plane
After=network.target nginx.service
Wants=nginx.service

[Service]
Type=simple
User=kokovoicebot
Group=kokovoicebot
WorkingDirectory=/opt/kokovoicebot
ExecStart=/opt/kokovoicebot/venv/bin/python oracle/app.py
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=3

# Environment variables — secrets
# IMPORTANT: Replace placeholders with actual values before enabling
Environment=TELEGRAM_BOT_TOKEN=%TELEGRAM_BOT_TOKEN%
Environment=ALLOWED_TELEGRAM_USER_ID=%ALLOWED_TELEGRAM_USER_ID%
Environment=TELEGRAM_WEBHOOK_SECRET=%TELEGRAM_WEBHOOK_SECRET%
Environment=GITHUB_REPO_OWNER=AiCurv
Environment=GITHUB_REPO_NAME=kokovoice
Environment=GITHUB_DISPATCH_TOKEN=%GITHUB_DISPATCH_TOKEN%
Environment=ORACLE_COMPLETION_SECRET=%ORACLE_COMPLETION_SECRET%
Environment=WEBHOOK_PORT=8443
Environment=WEBHOOK_URL_BASE=https://130.210.12.9
Environment=PYTHONPATH=/opt/kokovoicebot

# Security hardening
NoNewprivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/kokovoicebot/oracle/sessions.json /tmp
PrivateTmp=false

[Install]
WantedBy=multi-user.target
